# com.configurewise.leadmachine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.configurewise.leadmachine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/renko/docker/leadmachine/leadmachine-swift/.build/release/LeadMachineCLI</string>
        <string>process</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/renko/docker/leadmachine</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/renko/Library/Logs/LeadMachine/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/renko/Library/Logs/LeadMachine/stderr.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
